# daemon_test/system/d1_test.service
[Unit]
Description=Daemon1 test

[Service]
ExecStart=/tmp/daemon1.sh

[Install]
WantedBy=test.target


