# daemon_test/system/d3_test.service
[Unit]
Description=Daemon3 test
StopWhenUnneeded=yes

[Service]
ExecStart=/tmp/daemon3.sh

[Install]
WantedBy=test.target


