# daemon_test/system/test.target
[Unit]
Description=MyTest target
Requires=multi-user.target

[Install]
WantedBy=multi-user.target

