# wait_another_service/job1.target
[Unit]
Description=TestJob
DefaultDependencies=no
Conflicts=shutdown.target

